# queue_worker.ini
[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php ENVOY_ARTISAN queue:work --sleep=5 --tries=5
autostart=true
autorestart=true
user=root
numprocs=8
redirect_stderr=true
stdout_logfile=ENVOY_STORAGE/logs/worker.log 