# .dockerfiles/cron/crontab
* * * * * cd /var/www/html/ && php artisan schedule:run >> /dev/null 2>&1
# An empty line is required at the end of this file for a valid cron file.